# upgrade-files/rootfs/etc/crontabs/root
*/10   *  *   *   *     /usr/bin/cgminer-monitor
